# systemd/kismet@.service
[Unit]
Description=Kismet WiFi Scanner on interface %i
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/pi/AirscopeGuardian

Environment="PATH=/home/pi/AirscopeGuardian/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="DJANGO_SETTINGS_MODULE=config.settings"
Environment="PYTHONUNBUFFERED=1"

# Ensure log directory exists before starting
ExecStartPre=/bin/bash -c 'mkdir -p /home/pi/AirscopeGuardian/kismet/logs'

# Start Kismet via run_kismet.sh (foreground — systemd tracks the process)
ExecStart=/home/pi/AirscopeGuardian/kismet/root-shell/run_kismet.sh %i

# After kismet stops, import the latest .kismet DB
ExecStopPost=/home/pi/AirscopeGuardian/kismet/root-shell/import_kismet.sh

StandardOutput=append:/home/pi/AirscopeGuardian/kismet/logs/kismet_output.log
StandardError=append:/home/pi/AirscopeGuardian/kismet/logs/kismet_output.log
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
