# systemd/waveshare.service
# TEMPLATE — run install.sh to deploy. Do not copy manually.
[Unit]
Description=Waveshare display service for AirscopeGuardian
After=multi-user.target

[Service]
User=@@SERVICE_USER@@
Group=www-data
WorkingDirectory=@@PROJECT_DIR@@

Environment="PATH=@@PROJECT_DIR@@/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment=PYTHONUNBUFFERED=1

ExecStart=@@PROJECT_DIR@@/venv/bin/python @@PROJECT_DIR@@/gpio/waveshare.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
